# etc/systemd/system/snap-snapd-24792.mount
[Unit]
Description=Mount unit for snapd, revision 24792
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/snapd_24792.snap
Where=/snap/snapd/24792
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
